# r200_install/ds_connectivity_workaround/00-usb-R200.rules
ACTION=="add", ATTRS{idVendor}=="8086", ATTRS{idProduct}=="0a80", RUN+="/usr/local/bin/usb-R200-in_udev"
